# postgres-appliance/runit/etcd/run
#!/bin/sh -e

exec 2>&1
exec env -i ETCD_UNSUPPORTED_ARCH=$(dpkg --print-architecture) /bin/etcd --data-dir /run/etcd.data
