# postgres-appliance/runit/cron/run
#!/bin/sh -e

cmd="env -i supercronic $RW_DIR/cron/crontab/postgres"
exec 2>&1
if [ "$(id -u)" -eq 0 ]; then
    exec runuser -u postgres -- $cmd
fi
exec $cmd
