# postgres-appliance/runit/cron/run
#!/bin/sh -e

exec 2>&1
exec env -i /usr/sbin/cron -f
