# postgres-appliance/supervisor.d/cron.conf
[program:cron]
priority=100
directory=/
command=env -i /usr/sbin/crond -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
